# apps/player/systemd/misell-log-rotate.service
[Unit]
Description=misell-player log rotation

[Service]
Type=oneshot
Environment=MISELL_HOME=__MISELL_HOME__
EnvironmentFile=-%h/.config/misell-player/env
ExecStart=__MISELL_HOME__/scripts/rotate-logs.sh
